# debian/postinst
#!/bin/sh
sudo systemctl daemon-reload
sudo systemctl enable e-sabit.service
sudo systemctl start e-sabit.service