# debian/e-sabit/DEBIAN/postrm
#!/bin/sh
systemctl disable e-sabit.service
systemctl stop e-sabit.service
systemctl daemon-reload